# s6/services.d/2.tail/run
#!/bin/sh

exec tail -F /var/log/container.log -F /var/log/node/current
