# s6/services.d/1.node/run
#!/usr/bin/with-contenv sh

sleep 2

cd /app

su node

DEBUG=node.* npm start
